# over/etc/s6/filter/run
#!/usr/bin/execlineb -P
jq --unbuffered -crf data/script
