# over/etc/s6/addblock/up
#!/usr/bin/execline -P
ip a a 100.64.36.16/10 dev docker0
