# over/etc/s6/notify/run
#!/usr/bin/execlineb -P
chroot /mnt
curl -NsS --unix-socket /var/run/docker.sock 0/v1.20/events
